# deploy/vinaris-marketing-agent.service.example
[Unit]
Description=Vinaris Marketing Agent
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/vinaris-marketing-agent
Environment=HOST=127.0.0.1
Environment=PORT=8000
Environment=WORKERS=2
ExecStart=/opt/vinaris-marketing-agent/deploy/start_ubuntu.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
